# deploy/sales-heartbeat.timer
[Unit]
Description=Hourly heartbeat check for Sales Agent Chat Monitor
Requires=sales-heartbeat.service

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
