# deploy/sales-backup.timer
[Unit]
Description=Daily DB backup for Sales Agent Bot
Requires=sales-backup.service

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
